# dist/solartankstelle.service
[Unit]
Description=GUI für Solartankstelle

[Service]
Type=simple
User=1000
Group=1000
ExecStart=/usr/bin/python3 /opt/Solartankstelle_GUI/gui_main.py
WorkingDirectory=/opt/Solartankstelle_GUI
Environment=QT_QPA_PLATFORM=eglfs
Environment=QT_QPA_EGLFS_HIDECURSOR=1

[Install]
WantedBy=multi-user.target
